# service/com.blueorbit.examsgpt-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blueorbit.examsgpt-web</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/johelcastillo/Documents/BlueOrbit/examsGPT/service/examsgpt-web-wrapper.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/johelcastillo/Documents/BlueOrbit/examsGPT/server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>3000</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/johelcastillo/Documents/BlueOrbit/examsGPT/logs/web-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/johelcastillo/Documents/BlueOrbit/examsGPT/logs/web-server-error.log</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
